# debian/lernstick-config.btrfs_snapshots@.service
[Unit]
Description=Create Btrfs snapshot
After=lernstick-persistence.service

[Service]
ExecStart=create_btrfs_snapshot %i
