# debian/lernstick-config.lernstick-xserver-xorg-pre.service
[Unit]
Description=Updating the xorg configuration (pre)
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
Before=basic.target

[Service]
Type=oneshot
ExecStart=/lib/systemd/lernstick-xserver-xorg-pre

[Install]
WantedBy=basic.target
