# etc/systemd/system/sddm.service.d/override.conf
[Unit]
ConditionKernelCommandLine=|desktop=kde
ConditionKernelCommandLine=|desktop=enlightenment

[Install]
Alias=sddm.service
WantedBy=graphical.target
